# crontabs/mysql-monitor
# monitor mysql, restart if not running, then log actions with timestamp.
# run in root's crontab, or a user/service-account with sufficient privileges to restart mysql service
*/10 * * * * systemctl is-active --quiet mysql || (date +"%D %T mysql is down, restarting..." >> mysql-restart-log.txt && systemctl restart mysql)
